# systemd/root-protection-snapshot.timer
[Unit]
Description=Daily root-protection snapshot timer

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
